# samba/rootfs/etc/profile.d/startup.sh
#!/usr/bin/env bash

alias ll='ls -lh'
alias la='ls -lAh'